# board/canaan/k230-soc/rootfs_overlay/etc/profile.d/canaan_banner.sh
#!/bin/sh

cat <<EOF

$(echo -e "\033[1m")
    Welcome to K230, for more information, please visit:$(echo -e "\033[0m")
        https://www.kendryte.com/zh/document
        https://www.kendryte.com/answer/
        https://www.kendryte.com/zh/products
        https://github.com/kendryte/k230_linux_sdk/
        https://www.kendryte.com/k230_linux/dev/

EOF
